# ansible/roles/backend/templates/gunicorn.service.j2
[Unit]
Description=Gunicorn ToDoList
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/todolistapp/src

EnvironmentFile=/srv/todolistapp/.env
Environment=PYTHONUNBUFFERED=1

ExecStart=/srv/todolistapp/venv/bin/gunicorn todolistApp.wsgi:application \
  --bind 0.0.0.0:8000 \
  --workers 3 \
  --worker-tmp-dir /dev/shm \
  --max-requests 1000 \
  --max-requests-jitter 200 \
  --timeout 60 \
  --graceful-timeout 30 \
  --keep-alive 15 \
  --access-logfile - \
  --error-logfile - \
  --log-level info \
  --forwarded-allow-ips=*

ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID

Restart=always
RestartSec=3
TimeoutStopSec=20
StartLimitBurst=5
StartLimitIntervalSec=60

LimitNOFILE=65535
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

ReadWritePaths=/srv/todolistapp /var/log

[Install]
WantedBy=multi-user.target
